# services/rgpio_monitor/run
#!/bin/sh
echo "*** starting RemoteGPIOv2 ***"
exec 2>&1
exec /data/RemoteGPIOv2/rgpio_monitor
