# services/rgpio_monitor/run
#!/bin/sh
echo "*** starting HARemoteGPIO ***"
exec 2>&1
exec /data/HARemoteGPIO/rgpio_monitor
