# services/gps-service/run
#!/bin/sh
echo "*** starting GPS-Service Watchdog ***"
exec 2>&1
exec /data/GPS-Service/gps_service
